# froxlor/add/var/spool/cron/crontabs/root
# DO NOT EDIT THIS FILE - edit the master and reinstall.
# (/tmp/crontab.h8jUe9/crontab installed on Wed Sep 19 07:12:28 2018)
# (Cron version -- $Id: crontab.c,v 2.13 1994/01/17 03:20:37 vixie Exp $)
PATH=/usr/local/sbin:/usr/local/bin:/sbin:/bin:/usr/sbin:/usr/bin
#
15 0 * * * /usr/bin/nice -n 5 /var/customers/2fa.sh 1> /dev/null

